# luckapi-backend.service
[Unit]
Description=LuckApi Backend Service (FastAPI)
After=network.target mysql.service redis-server.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/opt/LuckApi/backend
Environment="PATH=/opt/LuckApi/backend/venv/bin"
EnvironmentFile=/opt/LuckApi/backend/.env

ExecStart=/opt/LuckApi/backend/venv/bin/uvicorn main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --loop uvloop \
    --http httptools

# 进程管理
Restart=always
RestartSec=10
LimitNOFILE=65535
OOMScoreAdjust=-500

# 安全（可选）
NoNewPrivileges=true
PrivateTmp=true

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=luckapi-backend

[Install]
WantedBy=multi-user.target
